# templates/default/jira.init.erb
#!/bin/sh
# JIRA startup script
#chkconfig: 2345 80 05
#description: JIRA
 
# Define some variables
# Name of app ( JIRA, Confluence, etc )
APP=jira
# Name of the user to run as
USER=jira
# Location of application's bin directory
BASE=<%= @jira_base %>
# Location of Java JDK
#export JAVA_HOME=/usr/lib/jvm/java-6-sun
# Location of Catalina PID
PIDFILE=$BASE/work/catalina.pid

. /etc/rc.d/init.d/functions

rh_status() {
    # run checks to determine if the service is running or use generic status
    echo $PIDFILE
    status -p $PIDFILE java
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}
 
case "$1" in
  start)
    echo "Starting $APP"
    /bin/su -m $USER -c "cd $BASE/logs && $BASE/bin/startup.sh &> /dev/null"
    echo "$APP started successfully"
    ;;
  stop)
    echo "Stopping $APP"
    /bin/su -m $USER -c "$BASE/bin/shutdown.sh &> /dev/null"
    echo "$APP stopped successfully"
    ;;
  restart)
    $0 stop
    sleep 5
    $0 start
    ;;
  status)
    rh_status
    ;;
  *)
    echo "Usage: /etc/init.d/$APP {start|restart|stop|status}"
    exit 1
    ;;
esac
 
exit 0
